# solitude.ini
[program:solitude]
command=python /data/solitude/app.py
numprocs=1
stopsignal=INT
directory=/data/solitude
user=root
redirect_stderr=true
stdout_logfile=/data/log/solitude/solitude.out
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stopasgroup=true
killasgroup=true
environment=PYTHONPATH=/data/solitude
